# deployment/crontab
# Laravel Scheduler
* * * * * cd /var/www/dartshop_laravel && php artisan schedule:run >> /dev/null 2>&1

# Backup database daily at 2:00 AM
0 2 * * * mysqldump -u dartshop -p'strong_db_password' dartshop > /var/backups/dartshop/backup-$(date +\%Y\%m\%d).sql

# Clean old logs weekly
0 3 * * 0 find /var/www/dartshop_laravel/storage/logs -name "*.log" -type f -mtime +30 -delete

# Clean old cache files daily
0 4 * * * cd /var/www/dartshop_laravel && php artisan cache:prune-stale-tags >> /dev/null 2>&1 